# whistle/whistle.ini
[group:whistle]
programs=whistled

[program:whistled]
command = w2 run
process_name=%(program_name)s
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0